# conf/varnish/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/varnishd
user = root
autostart = true
